# infra/hetzner/metastats-daily-crawl.service
[Unit]
Description=metastats TFT daily stats crawl (all regions, all ranks)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/metastats-crawler
EnvironmentFile=/etc/metastats-crawler/env
# 17 regions × ~30 min worst-case = ~9h. 20h ceiling matches the
# marketvalue unit so two daily timers always find a dead service.
TimeoutStartSec=20h
ExecStart=/usr/bin/node scripts/crawl-allranks-all-regions.mjs --mode auto
StandardOutput=journal
StandardError=journal
Nice=10
IOSchedulingClass=best-effort

[Install]
WantedBy=multi-user.target
